# .config/systemd/user/playerctld-daemon.service
[Unit]
Description=Playerctl Daemon

[Service]
Restart=always
RestartSec=1s

ExecStart=/usr/bin/playerctld daemon

[Install]
WantedBy=default.target
